# infrastructure/borgbackup.service
[Unit]
Description=Start Backup with Borg

[Service]
ExecStart=/etc/dojamdbbackup.sh
Type=simple
